# crontab
*/5	*	*	*	*	/bin/monitorscript/monitor.sh